# systemd/gravatar-avatar-sync.service
[Unit]
Description=Sync account avatar from Gravatar
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gravatar-avatar-sync
